# conf/supervisor/supervisord.conf
[program:nginx]
command=nginx

[program:php5-fpm]
command=php5-fpm

